# src/main/assembly/dist/bin/easy-bag-store.service
[Unit]
Description=EASY Bag Store Service

[Service]
ExecStart=/usr/local/easy-bag-store/bin/easy-bag-store-service-start

User=easy-bag-store
Group=easy-bag-store

[Install]
WantedBy=multi-user.target
